# coin-service-conf/node-exporter.service.erb
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=root
Group=root
ExecStart=/home/ubuntu/node-exporter/node_exporter
StandardOutput=append:/home/ubuntu/node-exporter/exporter.log
StandardError=append:/home/ubuntu/node-exporter/exporter-error.log
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=300
RestartSec=10s

[Install]
WantedBy=default.target
