# coin-service-conf/prometheus.service.erb
[Unit]
Description=Prometheus
After=network.target

[Service]
User=root
Group=root
ExecStart=/home/ubuntu/prometheus/prometheus --config.file=/home/ubuntu/prometheus/prometheus.yml
StandardOutput=append:/home/ubuntu/prometheus/prometheus.log
StandardError=append:/home/ubuntu/prometheus/prometheus-error.log
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=300
RestartSec=10s

[Install]
WantedBy=default.target
