# coin-service-conf/solana.service.erb
[Unit]
Description=Solana Daemon
After=network.target

[Service]
User=ubuntu
Group=ubuntu
LimitNOFILE=1000000
Environment=RUST_LOG=info
ExecStart=<%= start_command %>
StandardOutput=append:/home/ubuntu/solana-cluster/solana-validator.log
StandardError=append:/home/ubuntu/solana-cluster/solana-validator-error.log
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=300
RestartSec=10s

[Install]
WantedBy=default.target
