# coin-service-conf/solana-sys-tuner.service.erb
[Unit]
Description=Solana Sys Tuner Daemon
After=network.target

[Service]
User=root
Group=root
LimitNOFILE=1000000
ExecStart=<%= start_command %>
StandardOutput=append:/home/ubuntu/solana-cluster/solana-sys-tuner.log
StandardError=append:/home/ubuntu/solana-cluster/solana-sys-tuner-error.log
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=300
RestartSec=10s

[Install]
WantedBy=default.target
